# com.mohtimeos.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for Time OS API server.

  To install:
    cp com.mohtimeos.api.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.mohtimeos.api.plist

  Uses $HOME expansion via zsh -c (HOME is inherited from user session).
  Logs go to /tmp for simplicity.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mohtimeos.api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd "$HOME/clawd/moh_time_os" &amp;&amp; source .venv/bin/activate &amp;&amp; python api/server.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/mohtimeos-api.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mohtimeos-api.err.log</string>
</dict>
</plist>
