# com.mohtimeos.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for Time OS API server.

  Install:
    cp com.mohtimeos.api.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.mohtimeos.api.plist

  Uninstall:
    launchctl bootout gui/$(id -u)/com.mohtimeos.api

  Logs:
    tail -f ~/Library/Logs/mohtimeos-api.out.log
    tail -f ~/Library/Logs/mohtimeos-api.err.log
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mohtimeos.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/molhamhomsi/clawd/moh_time_os/.venv/bin/python3</string>
        <string>-m</string>
        <string>api.server</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/molhamhomsi/clawd/moh_time_os</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/molhamhomsi/clawd/moh_time_os</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/molhamhomsi/Library/Logs/mohtimeos-api.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/molhamhomsi/Library/Logs/mohtimeos-api.err.log</string>
</dict>
</plist>
